# scripts/systemd/riviere-forecast.service
[Unit]
Description=Riviere daily forecast

[Service]
Type=oneshot
User=riviere
WorkingDirectory=/home/riviere/riviere
ExecStart=/home/riviere/riviere/scripts/daily_forecast.sh
